# dynamic-layers/meta-virtualization/recipes-containers/docker/docker-moby/systemd/docker-image-prune.service
[Unit]
Description=Docker Image Prune

[Service]
Type=oneshot
ExecStart=/usr/bin/docker image prune --force
WorkingDirectory=/tmp
