# dynamic-layers/meta-virtualization/recipes-containers/docker/docker-moby/systemd/docker-container-prune.service
[Unit]
Description=Docker Container Prune

[Service]
Type=oneshot
ExecStart=/usr/bin/docker container prune --force
WorkingDirectory=/tmp
