# recipes-extended/klish/klish3/klishd.service
[Unit]
Type=simple
ExecStart=/usr/bin/klishd -d

[Install]
WantedBy=multi-user.target
