# recipes-extended/klish/klish2/konfd.service
[Unit]
Type=forking
PIDFile=/var/run/konfd.pid
ExecStart=/usr/bin/konfd

[Install]
WantedBy=multi-user.target
