# dynamic-layers/meta-virtualization/recipes-containers/podman/podman/systemd/podman-image-prune.service
[Unit]
Description=Podman Image Prune

[Service]
Type=oneshot
ExecStart=/usr/bin/podman image prune --force
WorkingDirectory=/tmp
