# dynamic-layers/meta-virtualization/recipes-containers/podman/podman/systemd/podman-container-prune.service
[Unit]
Description=Podman Container Prune

[Service]
Type=oneshot
ExecStart=/usr/bin/podman container prune --force
WorkingDirectory=/tmp
