# dynamic-layers/meta-virtualization/recipes-containers/podman/podman/systemd/podman-container-prune.timer
[Unit]
Description=Podman Container Prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
